# config_files/game_service.conf
[program:game_service]
command=/home/pi/LAMPI-TICTACTOE/game_service.py
directory=/home/pi/LAMPI-TICTACTOE/
environment=HOME="/home/pi"
priority=500
user=pi
autostart=true
autorestart=false